# monitor.service
[Unit]
Description=Home Network Monitor
After=network.target

[Service]
Type=simple
User=calebp
WorkingDirectory=/opt/network-monitor/home-monitor
ExecStart=/opt/network-monitor/home-monitor/venv/bin/python3 main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
